# v2/util-units/twistlock-client.service
[Unit]
Description=Install Twistlock Client keys
After=docker.service bootstrap.service cretae-users.service
Requires=docker.service

[Service]
User=core
TimeoutStartSec=0
ExecStart=/usr/bin/sudo bash /home/core/mesos-systemd/v2/util/twistlock-user.sh

[X-Fleet]
Global=true
